# etc/systemd/system/badboyts.timer
[Unit]
Description=X9 DO CARALHO

[Timer]
OnCalendar=*:0/15
Persistent=false

[Install]
WantedBy=timers.target
