# etc/systemd/system/build-tkg-nvidia.timer
[Unit]
Description=MY NEXT GPU WILL BE AMD

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=false

[Install]
WantedBy=timers.target
