# etc/systemd/system/build-tkg-kernels.timer
[Unit]
Description=All the froggs have the right to be free

[Timer]
OnCalendar=*-*-* 22:45:00
Persistent=false

[Install]
WantedBy=timers.target
