# etc/systemd/system/build-hourly.timer
[Unit]
Description=IT'S BUILD TIME

[Timer]
OnCalendar=*-*-* 00..03,05..18,20..23:00:00
Persistent=false

[Install]
WantedBy=timers.target
