# etc/systemd/system/build-morning.timer
[Unit]
Description=IT'S LLVM BUILD TIME

[Timer]
OnCalendar=*-*-* 4:10:00
Persistent=false

[Install]
WantedBy=timers.target
