# etc/systemd/system/build-night.timer
[Unit]
Description=IT'S FIREFOX BUILD TIME

[Timer]
OnCalendar=*-*-* 3:45:00
Persistent=false

[Install]
WantedBy=timers.target
