# etc/systemd/system/build-sunset.timer
[Unit]
Description=IT'S WINE BUILD TIME

[Timer]
OnCalendar=*-*-* 19:10:00
Persistent=false

[Install]
WantedBy=timers.target
